# systemd/osm-country-exports@daily.timer
[Unit]
Description=Daily OSM country exports sweep

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
